# zenodo-batch.service
[Unit]
Description=ICSAC Open Review Pipeline — twice-daily batch orchestrator (model check + watch tick + summary)
Documentation=https://github.com/3RiversWebTech/zenodo-pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=orangepi
Group=orangepi
WorkingDirectory=/home/orangepi/Desktop/icsac/zenodo-pipeline
EnvironmentFile=/home/orangepi/.config/zenodo-pipeline.env
ExecStart=/usr/bin/python3 /home/orangepi/Desktop/icsac/zenodo-pipeline/pipeline.py batch-tick
ExecStartPost=-/usr/bin/curl -fsS --max-time 10 "http://100.117.63.73:3001/api/push/bOaUZKHaJC?status=up&msg=OK&ping="
# Batch tick reviews all pending submissions in one shot with 3 passes each —
# worst case multiple papers × 3 passes × up to 5 slots + RQC. Extended timeout
# covers peak load.
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/orangepi/Desktop/icsac/zenodo-pipeline /home/orangepi/Desktop/icsac/icsacinstitute.org

[Install]
WantedBy=multi-user.target
