# zenodo-review.service
[Unit]
Description=ICSAC Zenodo AI-Assisted Review Pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=orangepi
WorkingDirectory=/home/orangepi/Desktop/icsac/zenodo-pipeline
ExecStart=/usr/bin/python3 /home/orangepi/Desktop/icsac/zenodo-pipeline/pipeline.py poll
EnvironmentFile=/home/orangepi/.config/zenodo-pipeline.env
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/orangepi/Desktop/icsac/zenodo-pipeline/reviews /home/orangepi/Desktop/icsac/zenodo-pipeline/downloads

[Install]
WantedBy=multi-user.target
